# docker/test/test-supervisord.conf
[supervisord]
nodaemon=true
loglevel=critical

[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf

[program:redis]
command=/usr/bin/redis-server

[program:login_frontend]
command=uwsgi -M --chdir /opt/app/sso_frontend --uwsgi-socket 127.0.0.1:3031 --module sso_frontend.wsgi --buffer-size 16384 --harakiri 5 --limit-nproc 40
stopsignal=INT

[program:p0f]
user=root
command=/opt/static/p0f-3.06b/p0f -i eth0 -f /opt/static/p0f-3.06b/p0f.fp -s /var/local/p0f/p0f.sock -u p0f "port 80 or port 443"

[program:websockets]
user=app
chdir=/opt/static/node_socket
command=/usr/bin/node /opt/static/node_socket/app.js